# examples/com.projectdiscovery.pdcp-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectdiscovery.pdcp-agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pdcp-agent</string>
        <string>-agent-output</string>
        <string>/Users/YOUR_USERNAME/.pdcp-agent/output</string>
        <string>-verbose</string>
        <string>-agent-tags</string>
        <string>production</string>
        <string>-agent-id</string>
        <string>unique-agent-id</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PDCP_API_KEY</key>
        <string>your-api-key</string>
        <key>PDCP_API_SERVER</key>
        <string>https://api.projectdiscovery.io</string>
        <key>PUNCH_HOLE_HOST</key>
        <string>proxy.projectdiscovery.io</string>
        <key>PUNCH_HOLE_HTTP_PORT</key>
        <string>8880</string>
        <key>PDCP_TEAM_ID</key>
        <string>your-team-id</string>
        <key>PROXY_URL</key>
        <string>http://127.0.0.1:8080</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.pdcp-agent/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.pdcp-agent/logs/stderr.log</string>
</dict>
</plist>

